# deploy/rebuild-mau.timer
[Unit]
Description=Run rebuild-mau.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
